# local-joshua.service
[Unit]
Description=Local Joshua AI Agent — Autonomous OSINT Analyst
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Joshua-AI
Environment=JOSHUA_OLLAMA_URL=http://localhost:11434
Environment=JOSHUA_MODEL=joshua:latest
Environment=BLACKBOARD_URL=http://localhost:9700
Environment=COURT_RECORDS_URL=http://localhost:9800
Environment=JOSHUA_VOICE_HOST=localhost
Environment=JOSHUA_VOICE_PORT=9876
Environment=JOSHUA_INFERENCE_DEVICE=cuda
Environment=JOSHUA_POLL_INTERVAL=10
Environment=JOSHUA_LOG_FILE=%h/Joshua-AI/joshua.log
ExecStart=%h/Joshua-AI/venv/bin/python %h/Joshua-AI/agent.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
